# wshop.service
[Unit]
Description=wshop

[Service]
ExecStart=/root/runwshop.sh

[Install]
WantedBy=multi-user.target
